# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=php-fpm

[program:nodejs]
command=node /var/www/server.js
